# deployments/salvobase.service
[Unit]
Description=Salvobase - MongoDB-compatible document database
Documentation=https://github.com/inder/salvobase
After=network.target

[Service]
Type=simple
User=salvobase
Group=salvobase

ExecStart=/usr/local/bin/salvobase \
    --port 27017 \
    --httpPort 27080 \
    --bind_ip 127.0.0.1 \
    --datadir /var/lib/salvobase \
    --logLevel info \
    --logFormat json

ExecReload=/bin/kill -HUP $MAINPID

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/salvobase
LimitNOFILE=65536
LimitNPROC=65536

# Restart policy
Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=salvobase

[Install]
WantedBy=multi-user.target
